# msearchd/msearchd.rc
#!/bin/ksh

daemon="/usr/local/sbin/msearchd"

. /etc/rc.d/rc.subr

rc_cmd $1
